# etc/supervisor/conf.d/keydb.conf
[program:keydb]
command=/opt/keydb/bin/keydb-server /etc/keydb/keydb.conf

user=keydb

stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
